# units/poller.timer
[Unit]
Description=Run date.service every 10 minutes

[Timer]
OnCalendar=*-*-* *:*:00
