# units/poller.timer
[Unit]
Description=Run poller.service every 1 minutes

[Timer]
OnCalendar=*-*-* *:*:00
